# data.mount
[Mount]
What=10.0.10.211:/nfs
Where=/data
Type=nfs
Options=noauto

[X-Fleet]
Global=true
