# systemd-units/srv-nginx.mount
[Mount]
What=10.0.10.211:/nfs/nginx
Where=/srv/nginx
Type=nfs
Options=noauto

[X-Fleet]
Global=true
